# src/main/rpm/postinstall.sh
chkconfig rpmtest on
service rpmtest start